# .config/systemd/user/habit.service
[Unit]
Description=Reminder daemon

[Service]
ExecStart=%h/.local/bin/habit.sh
Environment=DISPLAY=:0
