# .config/systemd/user/yt.service
[Unit]
Description=Download videos from yt everyday

[Service]
ExecStart=%h/.local/bin/yt.sh

