# .config/systemd/user/yt.timer
[Unit]
Description=Run yt service

[Timer]
OnCalendar=11:30:00

[Install]
WantedBy=multi-user.target
